# gcp-jobs/bn-retry/schedules/run.dev.cron
# Daily at 6AM Pacific Time (dev environment)
0 6 * * *
